# infra/open-webui/walter/open-webui-walter.service
[Unit]
Description=OvernightDesk Walter Open WebUI
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=simple
User=open-webui-walter
Group=open-webui-walter
SupplementaryGroups=docker
RuntimeDirectory=open-webui-walter
RuntimeDirectoryMode=0750
Environment=OPEN_WEBUI_IMAGE=ghcr.io/open-webui/open-webui@sha256:0d58a66704d69e52da83f72bcd43869ad4fd0c761313778bc95ef6940a0b81e3
Environment=OPEN_WEBUI_DATA_VOLUME=open-webui-hermes-walter-data
Environment=OPEN_WEBUI_PERSONA_CONFIG=/opt/open-webui-walter/source/persona-model.json
Environment=OPEN_WEBUI_PERSONA_SEED_SCRIPT=/opt/open-webui-walter/bin/seed_persona_model.py
ExecStartPre=+/opt/open-webui-walter/bin/load-phase-env.sh
ExecStartPre=+/opt/open-webui-walter/bin/seed-persona-model.sh
ExecStart=/opt/open-webui-walter/bin/run-container.sh
ExecStop=/opt/open-webui-walter/bin/stop-container.sh
Restart=on-failure
RestartSec=10
TimeoutStartSec=240
TimeoutStopSec=45
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/run/open-webui-walter

[Install]
WantedBy=multi-user.target
