# tenants/hermes-titus/email-poller/runtime/hermes-email-intake@.service
[Unit]
Description=OvernightDesk Hermes email intake (%i)
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=simple
User=titus-email-poller
Group=titus-email-poller
SupplementaryGroups=docker
ExecStartPre=+/opt/hermes-email-intake/bin/load-phase-config.sh %i
ExecStartPre=+/opt/hermes-email-intake/bin/prepare-volume.sh %i
ExecStart=/opt/hermes-email-intake/bin/run-container.sh %i
ExecStop=+/opt/hermes-email-intake/bin/stop-container.sh %i
Restart=always
RestartSec=10
TimeoutStartSec=300
TimeoutStopSec=45
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true

[Install]
WantedBy=multi-user.target
